# deployments/kind-rdma/soft-roce.service
# Copyright 2026 NVIDIA CORPORATION
# SPDX-License-Identifier: Apache-2.0
#
# Runs the baked Soft-RoCE setup once, at node boot, so the kind-rdma node
# image configures rdma_rxe itself — no `docker exec` from the demo's run.sh.
# Best-effort: the rdma_rxe kernel module comes from the host, so this no-ops
# where it is unavailable (e.g. macOS/Docker Desktop).
[Unit]
Description=Soft-RoCE (rdma_rxe) setup for the network-operator demo
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/setup-soft-roce.sh

[Install]
WantedBy=multi-user.target
